# Meldebian1/crontab
* * * * * cp -f /home/wordpress/escenario/logs/access.log /var/log/wordpress/
* * * * * cp -f /home/cowrie/cowrie/var/log/cowrie/* /var/log/cowrie/
