# Melfirewall/crontab
* * * * * /root/rsync-hp1.sh
* * * * * /root/rsync-hp2.sh
